# root/etc/s6-overlay/s6-rc.d/svc-mod-daapd-shairport-sync/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

until [[ ! -z $(ps -ae | grep avahi-damon) ]]; do
     sleep 1s
done

exec \
     /usr/bin/shairport-sync
