# root/etc/s6-overlay/s6-rc.d/init-mod-daapd-shairport-sync/up
/etc/s6-overlay/s6-rc.d/init-mod-daapd-shairport-sync/run
